# deploy/edudev-crm.service
[Unit]
Description=EduDev CRM backend
After=network.target postgresql.service

[Service]
Type=simple
User=edudev
Group=www-data
WorkingDirectory=/var/www/edudev/backend
EnvironmentFile=/var/www/edudev/backend/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
